# resources/install/debian/postrm
#!/bin/bash
# postrm script for jicofo

set -e

# Load debconf
. /usr/share/debconf/confmodule


case "$1" in
    purge)

        # Clear the debconf variable
        db_purge

        # stop the jicofo
        if [ -x /etc/init.d/jicofo ]; then
            invoke-rc.d jicofo stop
        fi

        # clear user and group
        # ensure nothing is running as jicofo user
        if ps -u jicofo h; then
            `ps -u jicofo -o pid h | xargs kill`
        fi
        # and then remove the user
        if getent passwd jicofo > /dev/null ; then
            deluser jicofo
        fi

        # clear logs
        if [ -d "/var/log/jicofo" ]; then
            rm -rf "/var/log/jicofo"
        fi

        rm -f /etc/jitsi/jicofo/config
    ;;

    remove)
        rm -rf /usr/share/jicofo/.java
        rm -rf /usr/share/jicofo/.sip-communicator
     ;;

    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

db_stop

exit 0
